# etc/cron.d/socialfeed
*/10 *   *   *   *   django ( . /home/django/envs/wagtail/bin/activate && cd /home/django/gilmoreplace_2022/ && ./manage.py update_social_feed && deactivate ) > /dev/null 2>&1
